# deploy/systemd/maw-tmux.service
[Unit]
Description=MAW dedicated tmux server
After=default.target

[Service]
Type=forking
# tmux exits when its last session closes unless exit-empty is off. Without
# this the server vanishes the moment start-server returns (no sessions
# yet), which is exactly what "active (exited)" with no PID means.
ExecStart=/usr/bin/tmux -L maw start-server \; set-option -s exit-empty off
ExecStop=/usr/bin/tmux -L maw kill-server
KillMode=mixed
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
